# recipes-containers/edgescale-eds/edgescale-eds-solution/edgescale.service
[Unit]
Description=edgescale

[Service]
User=root
EnvironmentFile=-/etc/default/edgescale
ExecStart=./usr/local/edgescale/bin/startup.sh

[Install]
WantedBy=multi-user.target
